# net-ftp/bftpd/files/bftpd.init
#!/sbin/openrc-run
command=/usr/sbin/bftpd
command_args="-d"
